# modes/ai/legacy-crucible12/systemd/crucible12@.service
[Unit]
Description=Crucible12 llama-server (%i preset)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/crucible12/setup
ExecStart=/opt/crucible12/setup/run-%i.sh
Restart=on-failure
RestartSec=5
User=crucible12
Group=crucible12

# The run-*.sh presets pass --mlock to pin the multi-GB model in RAM. systemd
# services default to RLIMIT_MEMLOCK=8MiB (DefaultLimitMEMLOCK), and because
# this unit runs as an unprivileged User= with no CAP_IPC_LOCK it can't raise
# that ceiling itself — so without this, --mlock silently degrades to a no-op
# (llama.cpp logs an mlock warning and runs unpinned). LimitMEMLOCK sets
# RLIMIT_MEMLOCK directly and is honored regardless of capabilities.
LimitMEMLOCK=infinity

# GPU/RAM headroom reservation while a preset is active is the mode-switcher's
# job (cgroup/slice priority via distro-modectl), not this unit — see DESIGN.md §5/§4.

[Install]
WantedBy=multi-user.target
